# deploy/msrelay.service
# =============================================================================
# msrelay.service - systemd unit for the relay server
#
# Install:
#   sudo cp msrelay.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now msrelay
# =============================================================================

[Unit]
Description=Mass Relay Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/msrelay --port 7744 --token REPLACE_WITH_YOUR_TOKEN
Restart=always
RestartSec=5
# Harden the service
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
